# triennale/terzo-anno/amministrazione-di-sistemi/template/esame/roles/service_configuration/files/app.service
[Unit]
Description=My Custom App Service
After=network.target

[Service]
ExecStart=/usr/local/bin/my_custom_app
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
